# etc/profile.d/aliases-safety.sh
#################################
#Use these to prevent data loss.#
#################################

alias chgrp='chgrp --preserve-root'
alias chmod='chmod --preserve-root'
alias chown='chown --preserve-root'
alias cp='cp -i'
alias ln='ln -i'
alias mv='mv -i'
alias rm='rm -i --preserve-root'
